# Build/0.3/etc/cron.monthly/EFA-Monthly-cron
#!/bin/bash
################################################################
# EFA Monthly cron
# Version 0.1-20121016
################################################################

# Check for update
/usr/local/sbin/EFA-Update -cron